# supervisord.conf
[supervisord]
nodaemon=true

[program:appium]
user=appium
command=/bin/bash -c "appium ${APPIUM_ARGS}"

[program:adb]
command=/bin/bash -c "${ANDROID_HOME}/platform-tools/adb start-server"